# src/rootfs/etc/services.d/thumbor/run
#!/bin/bash
set -e -u

exec s6-setuidgid root thumbor -c /etc/thumbor/thumbor.conf
